# contrib/init/bitcoinrd.conf
description "bitcoinR Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BITCOINRD_BIN="/usr/bin/bitcoinrd"
env BITCOINRD_USER="bitcoinr"
env BITCOINRD_GROUP="bitcoinr"
env BITCOINRD_PIDDIR="/var/run/bitcoinrd"
# upstart can't handle variables constructed with other variables
env BITCOINRD_PIDFILE="/var/run/bitcoinrd/bitcoinrd.pid"
env BITCOINRD_CONFIGFILE="/etc/bitcoinr/bitcoinr.conf"
env BITCOINRD_DATADIR="/var/lib/bitcoinrd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # bitcoinrd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BITCOINRD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitcoinrd."
        echo "The setting must appear in $BITCOINRD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"bitcoinR Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BITCOINRD_PIDDIR"
    chmod 0755 "$BITCOINRD_PIDDIR"
    chown $BITCOINRD_USER:$BITCOINRD_GROUP "$BITCOINRD_PIDDIR"
    chown $BITCOINRD_USER:$BITCOINRD_GROUP "$BITCOINRD_CONFIGFILE"
    chmod 0660 "$BITCOINRD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BITCOINRD_PIDFILE" \
    --chuid $BITCOINRD_USER:$BITCOINRD_GROUP \
    --exec "$BITCOINRD_BIN" \
    -- \
    -pid="$BITCOINRD_PIDFILE" \
    -conf="$BITCOINRD_CONFIGFILE" \
    -datadir="$BITCOINRD_DATADIR" \
    -disablewallet \
    -daemon

